# autolight_0.1alpha-1/usr/local/bin/autolight/autolight.service
[Unit]
Description=Webcam based auto brightness

[Service]
Restart=on-failure
WorkingDirectory=/usr/local/bin/autolight
EnvironmentFile=/usr/local/bin/autolight/config
ExecStart=/usr/local/bin/autolight/.venv/bin/python3 main.py

[Install]
WantedBy=default.target
